# scripts/com.user.auto-scroll-direction.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.auto-scroll-direction</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-c</string>
        <string>cd /Users/scott/Documents/Scott/macos_projects && ./scripts/auto_scroll_direction.sh --apply</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/tmp/auto-scroll-direction.err</string>
    <key>StandardOutPath</key>
    <string>/tmp/auto-scroll-direction.out</string>
</dict>
</plist>
